# 4ndr0tools/4ndr0service/systemd/user/env_maintenance.timer
[Unit]
Description=Daily 4ndr0service Environment Maintenance

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
