# 4ndr0tools/4ndr0service/systemd/user/env_maintenance.service
[Unit]
Description=Run 4ndr0service environment verification

[Service]
Type=oneshot
ExecStart=%h/.local/bin/4ndr0service/main.sh --fix

[Install]
WantedBy=default.target
