# 4ndr0tools/4ndr0service/systemd/user/env_maintenance.timer
[Unit]
Description=Daily 4ndr0service environment verification

[Timer]
OnCalendar=daily
Persistent=true
Unit=env_maintenance.service

[Install]
WantedBy=timers.target
