# Install/rsnapshot/sysd_service_file/rsnapshot-interval.timer
[Unit]
Description=rsnapshot hourly backup

[Timer]
OnCalendar=
Persistent=true
Unit=rsnapshot@interval.service

[Install]
WantedBy=timers.target
