# suckless/mem-police/systemd/mem-police.service
# /etc/systemd/system/mem-police.service
# Manages the mem-police daemon for system-wide reliability.

[Unit]
Description=Memory Policing Daemon
After=network.target

[Service]
Type=forking
PIDFile=/run/mem-police.pid
# This path assumes the compiled binary is in /usr/bin. Adjust if necessary.
ExecStart=/usr/bin/mem-police
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
