# suckless/mem-police/systemd/mem-police.service
# /etc/systemd/system/mem-police.service
# Manages the mem-police daemon using its canonical path.

[Unit]
Description=Memory Policing Daemon
After=network.target

[Service]
Type=forking
PIDFile=/var/run/mem-police.pid
ExecStart=/usr/bin/mem-police
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
