# install/rsnapshot/systemd/rsnapshot-daily.timer
[Unit]
Description=rsnapshot daily backup

[Timer]
OnCalendar=daily
Persistent=true
Unit=rsnapshot@daily.service

[Install]
WantedBy=timers.target
