# security/network/brave/systemd/brave.service
# ~/.config/systemd/user/brave.service
# Manages the Brave instance via its canonical script and applies Cgroup limits.

[Unit]
Description=Brave Web Browser with Cgroup Resource Controls
After=graphical-session.target

[Service]
# The %h specifier is the systemd-native, canonical way to refer to the user's home directory.
ExecStart=/usr/local/bin/bravebackgrounded

# ---- Cgroup Resource Control ----
MemoryMax=8G
MemoryHigh=6G
CPUWeight=50
CPUQuota=200%

[Install]
WantedBy=graphical-session.target
